# s2erunner/metricbeat/secrets/profile.d/env.sh
grep -qxF '10.99.15.100 metricd' /etc/hosts || echo '10.99.15.100 metricd' >> /etc/hosts
